# service/log/run
#!/bin/sh 
exec /usr/bin/multilog t /var/log/gps-server